# test/fake_api.conf
start on runlevel [2345]
stop on runlevel [!2345]

exec /usr/bin/python /test/fake_api.py
