# tests/files/fake_api.service
[Service]
Type=simple
ExecStart=/bin/python3 /etc/fake_api/fake_api.py